# .config/autostart/redshift-gtk.desktop
[Desktop Entry]
Icon=redshift
Name[ru]=Redshift
Name[pt_BR]=Redshift
Terminal=false
Name[pl]=Redshift
Comment[da]=Justeringsværktøj for farvetemperatur
Comment[ru]=Инструмент регулирования цветовой температуры
GenericName[da]=Justering af farvetemperatur
GenericName[hu]=Szín hőmérséklet beállítás
Name[tr]=Redshift
Comment=Color temperature adjustment tool
GenericName=Color temperature adjustment
Exec=redshift-gtk
Name[lt]=Redshift
GenericName[es]=Ajuste de la temperatura de color
Name[es]=Redshift
Comment[hu]=Színhőmérséklet beállító eszköz
GenericName[nl]=Kleurtemperatuur aanpassing
GenericName[ru]=Настройка цветовой температуры
Type=Application
GenericName[fr]=Réglage de la température de couleur
Name[de]=Redshift
Name[hu]=Redshift
Comment[cs]=Nástroj pro přizpůsobení barevné teploty
Comment[lt]=Spalvos temperatūros reguliavimo įrankis
Comment[nl]=Kleurtemperatuur aanpassingsprogramma
GenericName[zh_CN]=色温调节
Comment[es]=Herramienta para el ajuste de la temperatura de color
GenericName[de]=Farbtemperatur Anpassung
GenericName[tr]=Renk sıcaklığı ayarı
Comment[it]=Strumento per la regolazione della temperatura del colore
Comment[tr]=Renk sıcaklığı ayarlama aracı
X-GNOME-Autostart-enabled=true
GenericName[it]=Regolazione della temperatura del colore
Name[fr]=Redshift
Categories=Utility;
Comment[zh_CN]=色温调节工具
StartupNotify=true
Hidden=false
Comment[fr]=Outil de réglage de la température de couleur
Version=1.0
Name=Redshift
Name[it]=Redshift
GenericName[pl]=Dostosowanie temperatury barwowej
Comment[pl]=Narzędzie do dostosowywania temperatury barwowej
Name[zh_CN]=红移
GenericName[lt]=Spalvos temperatūros reguliavimas
GenericName[cs]=Přizpůsobení barevné teploty
Comment[de]=Farbtemperatur Anpassungswerkzeug

